# deploy/ipam.service
[Unit]
Description=IP Address Management (IPAM)
Documentation=https://github.com/davidmoore-rogers/ip-management
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=ipam
Group=ipam
WorkingDirectory=/opt/ipam
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/ipam/.env
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ipam
PrivateTmp=true

# Logging — stdout/stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ipam

[Install]
WantedBy=multi-user.target
