# deploy/polaris.service
[Unit]
Description=Polaris — IP Management Tool
Documentation=https://github.com/davidmoore-rogers/polaris
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=polaris
Group=polaris
WorkingDirectory=/opt/polaris
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/polaris/.env
Restart=on-failure
RestartSec=5

# Security hardening
AmbientCapabilities=CAP_NET_BIND_SERVICE
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/polaris
PrivateTmp=true

# Logging — stdout/stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polaris

[Install]
WantedBy=multi-user.target
