# deploy/shelob.service
[Unit]
Description=Shelob — IP Management Tool
Documentation=https://github.com/davidmoore-rogers/shelob
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=shelob
Group=shelob
WorkingDirectory=/opt/shelob
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/shelob/.env
Restart=on-failure
RestartSec=5

# Security hardening
AmbientCapabilities=CAP_NET_BIND_SERVICE
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/shelob
PrivateTmp=true

# Logging — stdout/stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shelob

[Install]
WantedBy=multi-user.target
